# services/vienna-runtime/deploy/vienna-console.service
[Unit]
Description=Vienna Console Backend
Documentation=https://docs.vienna.ai
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vienna
Group=vienna
WorkingDirectory=/opt/vienna/console/server

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/opt/vienna/console/server/.env.production

# Execution
ExecStart=/usr/bin/node dist/server.js

# Restart policy
Restart=always
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vienna-console

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vienna/workspace /opt/vienna/runtime
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%
TasksMax=512

[Install]
WantedBy=multi-user.target
